# device/remarkable-sync.service
[Unit]
Description=Sync reMarkable notebooks to aspirant-cell
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/root/sync-to-server.sh
TimeoutStartSec=600
